# systemd/system/docker.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd \
 --default-shm-size="1G" \
 --host=fd:// \
 --storage-driver=overlay2
LimitMEMLOCK=infinity
LimitSTACK=67108864
